# systemd/user/power-profile-tray.service
[Unit]
Description=MSI power-profile tray indicator
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/power-profile-tray
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
